# docker/setup/container-files/etc/supervisor.d/permissions.conf
[program:nginx-permissions]
command = sh bin/supervisor/nginx_permissions.sh
directory = /var/www
autostart=true
autorestart=false
startsecs=10
startretries=0
stopsignal=TERM
stopwaitsecs=5